# scrape.cron
0 */6 * * * python3 -m ~/scraper/app.Scrape